# coordinator/rpn-coord.service
[Unit]
Description=RPN Direct coordinator (fleet edition — WSS + STUN + TURN + REST)
After=network.target wg-quick@wgd0.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/rpn-coord

# Secrets — CI writes these from GitHub secrets (never commit actual values).
# /opt/rpn-coord/env format:
#   FLEET_TOKEN=<random — baked into Pi SD card images>
#   ADMIN_TOKEN=<random — for admin API>
#   TURN_SECRET=<random — for TURN credentials>
EnvironmentFile=/opt/rpn-coord/env

ExecStart=/opt/rpn-coord/rpn-coord \
    -http          :8089 \
    -stun          :3479 \
    -stun2         :3481 \
    -turn          :3480 \
    -public-ip     65.20.80.3 \
    -stun-public   65.20.80.3:3479 \
    -stun2-public  65.20.80.3:3481 \
    -relay-public  65.20.80.3:3480 \
    -turn-secret   ${TURN_SECRET} \
    -fleet-token   ${FLEET_TOKEN} \
    -admin-token   ${ADMIN_TOKEN} \
    -db            /opt/rpn-coord/rpn.db \
    -wg-iface      wgd0

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
